# deploy/compose/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── PostgreSQL ────────────────────────────────────────────────────────────────
[program:postgresql]
command=gosu postgres postgres -D /var/lib/postgresql/data
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=SIGINT

; ── FastAPI (uvicorn) ─────────────────────────────────────────────────────────
; priority=20: supervisord start de API altijd ná PostgreSQL. De API herprobeert
; de DB-verbinding via SQLAlchemy — startretries=10 dekt de opstartvertraging.
[program:api]
command=/app/api/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 3000 --workers 1
directory=/app/api
autostart=true
autorestart=true
priority=20
startretries=10
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
